# supervisord_fpm.ini
[program:php-fpm]
command = docker-php-entrypoint php-fpm -F
autostart=true
autorestart=true
startretries=3
stderr_logfile=/disk/log/supervisor.d/php_fpm.err.log
stdout_logfile=/disk/log/supervisor.d/php_fpm.out.log
user=root
